# docker-compose/example-laravel/schedule/supervisor-cron.conf
[program:cron]
command=/usr/local/bin/php /var/www/artisan queue:work
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/cron.log
